# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=service nginx start

[program:postfix]
command=/etc/init.d/postfix start


