# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=service nginx start

[program:php-fpm]
command=/etc/init.d/php7.0-fpm start

[program:postfix]
command=/etc/init.d/postfix start



